# deploy/polymarket-arb.service
[Unit]
Description=Polymarket arbitrage trading bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roy
WorkingDirectory=/home/roy/polymarket-arb
EnvironmentFile=-/home/roy/polymarket-arb/.env
Environment=HTTP_PROXY=http://127.0.0.1:7890
Environment=HTTPS_PROXY=http://127.0.0.1:7890
Environment=ALL_PROXY=socks5://127.0.0.1:7890
Environment=NO_PROXY=localhost,127.0.0.1,192.168.0.0/16
ExecStart=/home/roy/polymarket-arb/venv/bin/python /home/roy/polymarket-arb/main.py
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
